# supervisord.conf
[supervisord]
nodaemon=true

[program:myapp]
command=/app/main
autostart=true
autorestart=true
stderr_logfile=/var/log/myapp.err.log
stdout_logfile=/var/log/myapp.out.log
